# systemd/hueComposer.service
[Unit]
Description=Starts hueComposer on boot
After=network-online.target
#After=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=3
ExecStart=/usr/bin/python /root/hueComposer/hueComposer.py -c /root/imgs/config.yml -b /home/pi/imgs/

[Install]
WantedBy=multi-user.target
